# src/pi/service/moon-unit.service
[Unit]
Description=Moon Unit

[Service]
Type=simple
User=root
Environment=GAUGER_PORT=/dev/ttyS0
Environment=GPIO_ENABLED=1
ExecStart=/usr/local/bin/node /home/pi/git/moon-unit/index.js

[Install]
WantedBy=multi-user.target
